# Packs/pai-apple-ecosystem/launchd/com.pai.apple-ecosystem.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pai.apple-ecosystem</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/benjamin/.bun/bin/bun</string>
        <string>run</string>
        <string>/Users/benjamin/EscapeVelocity/PersonalAI/PAI/Packs/pai-apple-ecosystem/src/http-server.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/benjamin/EscapeVelocity/PersonalAI/PAI/Packs/pai-apple-ecosystem</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PORT</key>
        <string>8081</string>
        <key>TELOS_METRICS_PATH</key>
        <string>/Users/benjamin/EscapeVelocity/PersonalAI/PAI/Packs/pai-telos-metrics/data/metrics.jsonl</string>
        <key>TELOS_CONFIG_PATH</key>
        <string>/Users/benjamin/EscapeVelocity/PersonalAI/PAI/Packs/pai-telos-metrics/data/kpi-config.yaml</string>
        <key>PATH</key>
        <string>/Users/benjamin/.bun/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/benjamin/Library/Logs/pai-apple-ecosystem.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/benjamin/Library/Logs/pai-apple-ecosystem.error.log</string>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
</dict>
</plist>
